# .blackbox/4-scripts/com.siso.mcp-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.siso.mcp-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shaansisodia/DEV/SISO-ECOSYSTEM/SISO-INTERNAL/.blackbox/4-scripts/mcp-monitor-daemon.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/shaansisodia/.mcp-logs/mcp-monitor.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shaansisodia/.mcp-logs/mcp-monitor.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/shaansisodia/DEV/SISO-ECOSYSTEM/SISO-INTERNAL</string>
</dict>
</plist>
